# deploy/systemd/monitor4dingtalk.service
[Unit]
Description=Monitor4DingTalk Server Resource Monitor
Documentation=https://github.com/Jiang0977/monitor4dingtalk
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=monitor
Group=monitor
WorkingDirectory=/opt/monitor4dingtalk
ExecStart=/root/miniconda3/bin/python /opt/monitor4dingtalk/src/main.py
ExecReload=/bin/kill -HUP $MAINPID

# 重启策略
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# 资源限制
LimitNOFILE=65536
MemoryMax=256M

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/monitor4dingtalk /var/log/monitor4dingtalk

# 环境变量
Environment=PYTHONPATH=/opt/monitor4dingtalk
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target 